# deployment/nemoclaw-blueprint/supervisord-sandbox.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log

[program:bus-agent]
command=/app/entrypoint-sandbox.sh
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
